# docker/conf/supervisor/cron.ini
[program:cron]
process_name=%(program_name)s_%(process_num)02d
command=crond -f
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopwaitsecs=10